# deploy/slidescanner.service
[Unit]
Description=Slide Scanner capture web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
SupplementaryGroups=plugdev
WorkingDirectory=__REPO__
ExecStart=/usr/bin/python3 __REPO__/capture_server.py --port __PORT__ --out-dir __OUT__ --no-setup
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
